# k8s/edge/li-httpd-homelab.service
[Unit]
Description=li-httpd homelab edge (LIS ingress)
After=network-online.target k3s.service
Wants=network-online.target

[Service]
Type=simple
Environment=LIC_ROOT=/home/s4il0r/staging/lic
Environment=LI_HTTPD_ROOT=/home/s4il0r/staging/li-httpd
Environment=MAJICO_HTTPD_TOML=/home/s4il0r/staging/majico-deploy/deploy/staging/edge/majico-staging.httpd.toml
ExecStartPre=/bin/mkdir -p /run/li-httpd /var/lib/li-httpd/empty
ExecStartPre=/bin/bash -c 'HOME=/home/s4il0r LIC_ROOT=/home/s4il0r/staging/lic LI_HTTPD_ROOT=/home/s4il0r/staging/li-httpd MAJICO_HTTPD_TOML=/home/s4il0r/staging/majico-deploy/deploy/staging/edge/majico-staging.httpd.toml /home/s4il0r/staging/beelink-cleanup/scripts/edge-lis-apply.sh --render-only'
ExecStart=/usr/local/bin/li-httpd /run/li-httpd/homelab.runtime.conf
User=root
Group=root
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
