# k8s/edge/li-httpd-edge-watchdog.timer
[Unit]
Description=Homelab li-httpd edge watchdog (every 5 minutes)

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
