# k8s/edge/li-httpd-homelab-tls.service
[Unit]
Description=li-httpd homelab edge TLS (:443)
After=network-online.target k3s.service li-httpd-homelab.service
Wants=network-online.target
# Render: edge-lis-apply.sh holds /run/li-httpd/edge-apply.lock (safe parallel restarts).

[Service]
Type=simple
Environment=LIC_ROOT=/home/s4il0r/staging/lic
Environment=LI_HTTPD_ROOT=/home/s4il0r/staging/li-httpd
Environment=MAJICO_HTTPD_TOML=/home/s4il0r/staging/majico-deploy/deploy/staging/edge/majico-staging.httpd.toml
ExecStartPre=/bin/mkdir -p /run/li-httpd /var/lib/li-httpd/tls/homelab /var/lib/li-httpd/empty
ExecStartPre=/bin/bash -c 'for i in $(seq 1 120); do test -f /run/li-httpd/.render-ready && exit 0; sleep 1; done; echo "timeout waiting for /run/li-httpd/.render-ready (HTTP unit must render first)" >&2; exit 1'
ExecStart=/usr/local/bin/li-httpd /run/li-httpd/homelab.tls.runtime.conf
ExecStartPost=-/usr/local/bin/edge-health-probe.sh 443 gitlab.lilangverse.xyz /health
User=root
Group=root
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
