# k8s/edge/li-httpd-edge-watchdog.service
[Unit]
Description=Homelab li-httpd edge watchdog (GitLab local HTTPS probe (WAN logged only))
After=network-online.target li-httpd-homelab.service li-httpd-homelab-tls.service

[Service]
Type=oneshot
Environment=HOME=/home/s4il0r
ExecStart=/bin/bash /home/s4il0r/staging/homelab-k3s/scripts/edge-watchdog.sh